# crontab
# crontab
*/5 * * * * /app/entrypoint.sh generate_jcasc >> /var/log/cron.log 2>&1
